# sysutils/logstash/files/FreeBSD.rc.d.logstash-receive
#!/bin/sh

# Configuration settings for logstash-receive in /etc/rc.conf:
#
# PROVIDE: logstash-receive
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown
#
# logstash_receive_enable (bool):
#   Default value: "NO"
#   Flag that determines whether Logstash is enabled.
#
# logstash_receive_home (string):
#   Default value: "/usr/local/logstash"
#   Logstash installation directory.
#
# logstash_receive_config (string):
#   Default value: /usr/local/etc/${name}/${name}.conf
#   Logstash configuration file path.
#
# logstash_receive_log (bool):
#   Set to "NO" by default.
#   Set it to "YES" to enable logstash logging to file
#   Default output to /var/log/logstash_receive.log
#
# logstash_receive_log_file (string):
#   Default value: "${logdir}/${name}.log"
#   Log file path.
#
# logstash_receive_java_home (string):
#   Default value: "/usr/local/openjdk8"
#   Root directory of the desired Java SDK.
#   The JAVA_HOME environment variable is set with the contents of this
#   variable.
#
# logstash_receive_java_opts (string):
#   Default value: ""
#   Options to pass to the Java Virtual Machine.
#   The JAVA_OPTS environment variable is set with the contents of this
#   variable.
#
# logstash_receive_opts (string):
#   Default value: ""
#   Additional command line flags for logstash, eg. "-r"
#

. /etc/rc.subr

name=logstash_receive
dashname=logstash-receive
rcvar=logstash_receive_enable

load_rc_config ${name}

logdir="/var/log"

: ${logstash_receive_enable="NO"}
: ${logstash_receive_home="/usr/local/logstash"}
: ${logstash_receive_config="/usr/local/etc/logstash/${dashname}.conf"}
: ${logstash_receive_log="NO"}
: ${logstash_receive_log_file="${logdir}/${dashname}.log"}
: ${logstash_receive_java_home="/usr/local/openjdk8"}
: ${logstash_receive_java_opts=""}
: ${logstash_receive_opts=""}

piddir=/var/run/logstash
pidfile=${piddir}/${dashname}.pid

if [ -d $piddir ]; then
  mkdir -p $piddir
fi

logstash_receive_cmd="${logstash_receive_home}/bin/logstash"
procname="${logstash_receive_java_home}/bin/java"

logstash_receive_chdir=${logstash_receive_home}
logstash_receive_log_options=""

if checkyesno logstash_receive_log; then
  logstash_receive_log_options=" --log ${logstash_receive_log_file}" 
fi

logstash_receive_args="agent -f ${logstash_receive_config} ${logstash_receive_log_options} ${logstash_receive_opts}"

JAVA_OPTS="${logstash_receive_java_opts}"
JAVA_HOME="${logstash_receive_java_home}"
export JAVA_OPTS
export JAVA_HOME

command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} ${logstash_receive_cmd} ${logstash_receive_args}"
required_files="${logstash_receive_home} ${logstash_receive_java_home} ${logstash_receive_cmd} ${logstash_receive_config}"
#echo command_args=$command_args
# Include /usr/local/bin in path because Logstash startup scripts
# assume bash is in path.
PATH=/usr/local/bin:$PATH

run_rc_command "$1"
